# listen.service
[Unit]
Description=Cross server command TCP listener

[Service]
Type=forking
ExecStart=/etc/cross-server-commands/listen
